# systemd/chrome-manager.service
[Unit]
Description=Chrome Isolation Manager Web Interface
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=sam
WorkingDirectory=/home/sam/.local/share/chrome-isolation-manager/app
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /home/sam/.local/share/chrome-isolation-manager/app/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
